# src/configs/local-macosx-m2-work/Library/LaunchAgents/com.workflow.weekly-notebook.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.workflow.weekly-notebook</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/<USER>/.claude/routines/weekly-notebook.sh</string>
	</array>

	<key>StartCalendarInterval</key>
	<dict>
		<key>Weekday</key>
		<integer>5</integer>
		<key>Hour</key>
		<integer>15</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>

	<key>StandardOutPath</key>
	<string>/Users/<USER>/.claude/routines/logs/weekly-notebook.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/<USER>/.claude/routines/logs/weekly-notebook.log</string>
</dict>
</plist>
