# src/configs/local-macosx-m2-work/Library/LaunchAgents/com.workflow.trello-daybook-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.workflow.trello-daybook-sync</string>

	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/cron/run_once_daily.sh</string>
		<string>/bin/bash</string>
		<string>/Users/<USER>/.claude/routines/trello-daybook-sync.sh</string>
	</array>

	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Weekday</key>
			<integer>1</integer>
			<key>Hour</key>
			<integer>7</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>2</integer>
			<key>Hour</key>
			<integer>7</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>3</integer>
			<key>Hour</key>
			<integer>7</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>4</integer>
			<key>Hour</key>
			<integer>7</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>5</integer>
			<key>Hour</key>
			<integer>7</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
	</array>

	<key>StandardOutPath</key>
	<string>/Users/<USER>/.claude/routines/logs/trello-daybook-sync.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/<USER>/.claude/routines/logs/trello-daybook-sync.log</string>
</dict>
</plist>
