# src/configs/local-macosx-m4/Library/LaunchAgents/com.user.mountsmb.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.mountsmb</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/osascript</string>
        <string>-e</string>
        <string>mount volume "smb://sthinds@sthinds.local/sthinds"</string>
        <string>mount volume "smb://sthinds@sthinds.local/D"</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>300</integer> <!-- retry every 5 mins -->
</dict>
</plist>
