# src/configs/macosx/Library/LaunchAgents/com.workflow.bucket-files-by-date.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.workflow.bucket-files-by-date</string>

	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/cron/bucket_files_by_date_targets.sh</string>
	</array>

	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Weekday</key>
			<integer>1</integer>
			<key>Hour</key>
			<integer>1</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>4</integer>
			<key>Hour</key>
			<integer>1</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
	</array>

	<key>StandardOutPath</key>
	<string>/Users/<USER>/.cache/.workflow/launchd.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/<USER>/.cache/.workflow/launchd.log</string>
</dict>
</plist>
