# src/configs/macosx/Library/LaunchAgents/com.user.host-relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.user.host-relay</string>

	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python3</string>
		<string>/usr/local/bin/host_relay/host_relay.py</string>
	</array>

	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>

	<key>StandardOutPath</key>
	<string>/tmp/host-relay.stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/host-relay.stderr.log</string>
</dict>
</plist>
